# conary-proxy/conary-proxy.cron
#!%(essentialbindir)s/bash
%(sbindir)s/tmpwatch 720 %(proxydir)s/cscache
%(sbindir)s/tmpwatch 720 %(proxydir)s/contents

