# slocate/slocate.cron
#!/bin/sh
. /etc/updatedb.conf
[ "$DAILY_UPDATE" != "yes" ] && exit 0

renice +19 -p $$ >/dev/null 2>&1

if [ -x /usr/bin/ionice ];then
    ionice -c3 /usr/bin/updatedb
else
    /usr/bin/updatedb
fi
