# oversite/oversite.cron
#!/bin/bash
#
# Script used to get Foresight updates
#

# Choose a random sleep interval so we don't beat up the server
sleep $(expr $RANDOM \% 720)

/usr/sbin/oversite.sh
