# popularity/popularity.cron
#!/bin/bash 
# Script used to show distrowatch how popular we are
#

# Choose a random (up to 12 hrs) sleep interval so we don't flood distrowatch
sleep $(expr $RANDOM \% 43200)


wget --user-agent="Mozilla/5.0 (X11; U; Linux i686; en-US; rv:1.7.6) Gecko/20050405 Firefox/1.0.2" --directory-prefix=/$HOME/.popularity http://distrowatch.com/foresight
