# bootsplash/init.bootsplashstop
#!/bin/bash
#
# bootsplashstop       Stop the bootsplash throbber
#
# chkconfig: - 99 99
# description: Stop the bootsplash throbber

. /etc/init.d/functions


# See how we were called.
case "$1" in
  start)
        echo $"Stopping bootsplash throbber: "
	/sbin/bootanim stop
        ;;
  *)
        echo $"Usage: $0 {start}"
        exit 1
esac

exit 0
