# bootsplash/init.bootsplash
#!/bin/bash
#
# bootsplash       Start the bootsplash throbber
#
# chkconfig: 5 01 99
# description: Start the bootsplash throbber

. /etc/init.d/functions


# See how we were called.
case "$1" in
  start)
        # Assign devices
        echo $"Starting bootsplash throbber: "
	/sbin/bootanim start throbber.mng
        ;;
esac

exit 0
